# systemd/hoistwaywatch-capture.service
[Unit]
Description=HoistwayWatch capture health service
After=network-online.target hoistwaywatch-nats.service
Wants=network-online.target hoistwaywatch-nats.service

[Service]
Type=simple
WorkingDirectory=/opt/hoistwaywatch
Environment=HW_NATS_URL=nats://127.0.0.1:4222
Environment=HW_SITE_ID=site_default
Environment=HW_CAMERA_ID=cam1
Environment=HW_CAMERA_SOURCE=0
ExecStart=/opt/hoistwaywatch/.venv/bin/hw-capture
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target

